# mangdang/System/rc.local
#!/bin/bash
#
# 1.Initialize power switch and pwm units
#
# Copyright (c) 2020, Mangdang Technology Co., Limited
#
#

#export pwm channel
for i in $(seq 0 15)
do
	echo $i > /sys/class/pwm/pwmchip0/export
	echo  4000000 > /sys/class/pwm/pwmchip0/pwm$i/period
	# chmod  777  /sys/class/pwm/pwmchip0/pwm$i/duty_cycle
done
mpg123 -q /home/ubuntu/Music/power_on.mp3 &