# mangdang/FuelGauge/battery_monitor.service
[Unit]
Description=max17205 battery monitor service

[Service]
Type=simple
ExecStart=/usr/bin/battery_monitor
User=root

[Install]
WantedBy=sysinit.target
